# install/autostart/skippy-xd.desktop
[Desktop Entry]
Type=Application
Exec=skippy-xd --start-daemon &
Icon=preferences-system-windows-move
Name=Skippy XD
