# install/autostart/touchegg.desktop
[Desktop Entry]
Type=Application
Name=Touchégg
Exec=touchegg
Icon=libinput-gestures
Comment=Touchégg Multi-Touch Gesture Client
